# files/21_linux_home/.config/systemd/user/bitbake.slice
[Unit]
Description=Slice for running Bitbake

[Slice]
MemoryHigh=208G
MemoryMax=224G
MemorySwapMax=0
