# files/21_linux_home/.config/systemd/user/bitbake.slice
[Unit]
Description=Slice for running Bitbake

[Slice]
MemoryHigh=32G
MemoryMax=40G
MemorySwapMax=0
