# files/21_linux_home/.config/systemd/user/bitbake.slice
[Unit]
Description=Slice for running Bitbake

[Slice]
MemoryHigh=48G
MemoryMax=52G
MemorySwapMax=0
