# files/21_linux_home/.config/systemd/user/obmc-jobber-nightly.service
[Unit]
Description=obmc-jobber Nightly

[Service]
Type=oneshot
Environment=SSH_AUTH_SOCK=/run/user/1000/keyring/ssh
ExecStart=/home/apwillia/.zinit/plugins/williamspatrick---dotfiles/bin/20_all/obmc-jobber -s nightly
