# files/21_linux_home/.config/systemd/user/offlineimap.timer
[Unit]
Description=offlineimap Timer

[Timer]
OnBootSec=2m
OnUnitInactiveSec=5m

[Install]
WantedBy=default.target
