# files/21_linux_home/.config/systemd/user/offlineimap.service
[Unit]
Description=offlineimap sync

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap
TimeoutStopSec=5m

[Install]
WantedBy=default.target
