# files/21_linux_home/.config/systemd/user/prefect-worker.service
[Unit]
Description=Prefect Worker Service (linux-dev)

[Service]
Type=simple
RuntimeDirectory=prefect-worker
WorkingDirectory=%t/prefect-worker
ExecStart=/bin/zsh -lc 'source ~/.zshrc && /usr/bin/uv run --project /home/apwillia/private/dev/prefect-automation prefect worker start --pool linux-dev --name desktop --limit 16 --work-queue default'
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prefect-worker

[Install]
WantedBy=default.target
