# files/99_facebook_linux/.config/systemd/user/unison.service
[Unit]
Description=Unison synchronization daemon

[Service]
ExecStart=unison -socket 9999

[Install]
WantedBy=default.target

# vi: ft=dosini
